# common/etc/crontab
MAILTO=iapp@example.org
15 * * * *  root  cd / && run-parts --report /etc/cron.hourly
0  0 * * *  root  test -x /usr/sbin/anacron || ( cd / && run-parts --report /etc/cron.daily )
